# systemd/admute.service
[Unit]
Description=Admute Service

[Service]
Type=simple
ExecStart=%h/.scripts/admute
Restart=on-failure

[Install]
WantedBy=multi-user.target
